# media-tv/plex-media-server/files/pms
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need localmount
	after net
	keyword -timeout
}

start() {
	ebegin "Starting plex-media-server"
	. /etc/default/plexmediaserver
	start-stop-daemon --start -b -x /usr/sbin/start_pms -u $PLEX_MEDIA_SERVER_USER --pidfile /var/run/pms.pid
	eend $?
}

stop() {
	ebegin "Stopping plex-media-server"
	start-stop-daemon --stop -x /usr/sbin/start_pms --pidfile /var/run/pms.pid
	eend $?
}
