# app-misc/qdrant/files/qdrant.service
# Copyright 1999-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Qdrant vector database
After=network.target

[Service]
Type=simple
User=qdrant
Group=qdrant
WorkingDirectory=/var/lib/qdrant
ExecStart=/usr/bin/qdrant --config-path /etc/qdrant/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
