# sys-apps/systemd-units/files/services-desktop/cups.service
[Unit]
Description=CUPS Printing Service

[Service]
ExecStart=/usr/sbin/cupsd -f

[Install]
WantedBy=multi-user.target
