# sys-apps/systemd-units/files/services-desktop/powernowd.service
[Unit]
Description=PowerNow! Daemon

[Service]
Type=simple
ExecStart=/usr/sbin/powernowd -d -q

[Install]
WantedBy=multi-user.target
